# service/jmr.service
[Unit]
Description=Jellyfin Media Renamer

[Service]
ExecStart=/usr/local/bin/jmr
User=jellyfin-media-renamer

[Install]
WantedBy=multi-user.target

