# service/jmr.service
[Unit]
Description=Jellyfin Media Renamer

[Service]
User=jellyfin-media-renamer
WorkingDirectory=$HOME/jmr-app
ExecStart=$HOME/jmr-app/start-all.sh

[Install]
WantedBy=multi-user.target

